# ansible/templates/app.service.j2
[Unit]
Description=app
After=network.target

[Service]
User=app
ExecStart=/usr/bin/java -jar /opt/app/app.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
